# tmux-shutdown-save.service
[Unit]
Description=Save tmux session before shutdown/logout

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
ExecStop=%h/.config/tmux/scripts/save-session-before-shutdown.sh
TimeoutStopSec=5

[Install]
WantedBy=default.target
